# DEBIAN/postinst
#!/bin/bash
systemctl enable key-mapper
systemctl stop key-mapper
systemctl start key-mapper
