# DEBIAN/postinst
#!/bin/bash
if [ -d "/run/systemd/system/" ]; then
  systemctl enable key-mapper
  systemctl stop key-mapper
  systemctl start key-mapper
fi
